# mini-cloud/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/apachectl -DFOREGROUND
autorestart=true

[program:nextcloud-cron]
command=php -f /var/www/html/cron.php
user=www-data
autorestart=true

[program:jupyter]
command=jupyter notebook --ip=0.0.0.0 --port=8888 --NotebookApp.token='' --NotebookApp.password='' --allow-root --no-browser --notebook-dir=/data/notebooks
autorestart=true

[program:prometheus]
directory=/opt/prometheus
command=/opt/prometheus/prometheus --config.file=/etc/mini-cloud/monitoring/prometheus.yml --storage.tsdb.path=/data/prometheus
autorestart=true

[program:grafana]
directory=/opt/grafana/bin
command=/opt/grafana/bin/grafana-server --homepath=/opt/grafana
autorestart=true

[program:samba]
command=/usr/sbin/smbd --foreground --log-stdout
autorestart=true

[program:website]
command=python3 -m http.server 8080 --directory /srv/website
autorestart=true

[program:chat-ui]
command=python3 -m http.server 3000 --directory /srv/chat-ui
autorestart=true

[program:llm-ui]
directory=/srv/llm-chat-ui
command=python3 app.py --host 0.0.0.0 --port 3001
autorestart=true